# debian/puppet-master.cron.daily
#!/bin/sh

if [ -e /var/lib/puppet/reports ] ; then
	find /var/lib/puppet/reports -maxdepth 2 -type f -ctime +30 -delete
fi

exit 0
